# untangle-lcd/debian/untangle-lcd.postinst
#!/bin/bash

ourInit() {
    if [ -x /usr/sbin/invoke-rc.d ] && [ ! "`readlink /usr/sbin/invoke-rc.d`" = /bin/true ] ; then
        invoke-rc.d $1 $2
    else
        /etc/init.d/$1 $2
    fi
}

if [ -x "/etc/init.d/untangle-lcd" ]; then
    update-rc.d untangle-lcd defaults 40 > /dev/null
    ourInit untangle-lcd restart
fi

exit 0
