# untangle-portwell-lcd/debian/untangle-portwell-lcd.service
[Unit]
Description=Untangle Portwell LCD

[Service]
Type=simple
User=daemon
ExecStart=/usr/sbin/eizo
Restart=always

[Install]
WantedBy=multi-user.target
